# service.d/01wm/run
#!/bin/bash
# -viewpasswd string  - need -passwd specified 
# -passwdfile  
#  -usepw
#  -storepasswd
#
#
sv start 00xvfb || exit 1
exec ratpoison 
